# init/goliath-0.conf
start on started haproxy
stop on stopping haproxy
respawn

chdir /app/src

exec bundle exec ruby hello.rb --socket /tmp/hello_socket --stdout >&5

