# init/goliath-3.conf
start on started haproxy
stop on stopping haproxy
respawn

chdir /app/src

exec bundle exec ruby kill.rb --socket /tmp/kill_socket --stdout >&5
