# init/goliath-2.conf
start on started haproxy
stop on stopping haproxy
respawn

chdir /app/src

exec bundle exec ruby ps.rb --socket /tmp/ps_socket --stdout >&5
