# init/goliath-1.conf
start on started haproxy
stop on stopping haproxy
respawn

chdir /app/src

exec bundle exec ruby error.rb --socket /tmp/error_socket --stdout >&5
