# init/haproxy.conf
start on startup
stop on session-end
respawn

chdir /app/src


exec haproxy -f /app/src/haproxy.cfg >&5